# systemd/mnt-floppy.mount
[Mount]
What=/dev/sda
Where=/mnt/floppy
